# alarmdecoder/supervisord.conf
[supervisord]
nodaemon=true

[program:services]
command=/bin/bash /usr/local/bin/start_services.sh
